# systemd/user/waybar.service.d/override.conf
[Service]
Environment="PATH=%h/tools/scripts/hypr-helpers:%h/go/bin:%h/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/mdeaguiargomes/.local/bin/"
